# services/backup/systemd/alona-core-backup.timer
[Unit]
Description=Run Alona Core Backup Daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target

